# roles/srsilo/templates/srsilo-update.service.j2
[Unit]
Description=srSILO Update Pipeline - Automated Genomic Data Processing
After=network-online.target docker.service
Wants=network-online.target
Requires=docker.service

[Service]
Type=oneshot
User={{ srsilo_user }}
WorkingDirectory={{ wisepulse_repo_path }}/srsilo-updater

ExecStart={{ srsilo_base_path }}/venv/bin/python -m pipeline \
    --config {{ srsilo_base_path }}/pipeline.yml

StandardOutput=journal
StandardError=journal
SyslogIdentifier=srsilo-update

Restart=no

PrivateTmp=yes
NoNewPrivileges=yes
ProtectSystem=strict
ReadWritePaths={{ srsilo_base_path }}

[Install]
WantedBy=multi-user.target
