# roles/srsilo/templates/wisepulse-pipeline.service.j2
[Unit]
Description=WisePulse Data Pipeline - Smart Fetch and Process
Documentation=https://github.com/cbg-ethz/WisePulse
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=oneshot
User={{ wisepulse_user }}
Group={{ wisepulse_group }}
WorkingDirectory={{ wisepulse_repo_path }}

# Environment variables for the pipeline
Environment="FETCH_API_BASE_URL={{ wisepulse_api_base_url }}"
Environment="FETCH_DAYS={{ wisepulse_fetch_days }}"
Environment="FETCH_MAX_READS={{ wisepulse_fetch_max_reads }}"
Environment="LAPIS_PORT={{ wisepulse_lapis_port }}"
Environment="RETENTION_DAYS={{ wisepulse_retention_days }}"
Environment="RETENTION_MIN_KEEP={{ wisepulse_retention_min_keep }}"
Environment="PATH=/home/{{ wisepulse_user }}/.cargo/bin:/usr/local/bin:/usr/bin:/bin"

# Run the smart pipeline (checks for new data first)
ExecStart=/usr/bin/make smart-fetch-and-process

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=wisepulse-pipeline

# Resource limits (optional, adjust as needed)
# MemoryMax=4G
# CPUQuota=200%

# Security settings
NoNewPrivileges=true
PrivateTmp=true

# Restart policy - don't restart on failure for scheduled jobs
Restart=no

[Install]
WantedBy=multi-user.target
